# share/system/install.d/restart.cron
#!/bin/bash

if [ -x /etc/init.d/cron ]; then
     /etc/init.d/cron restart
fi
