# modules/videowall/node/post-install.d/lib/udev/rules.d/65-xorg-nvidia.rules
ACTION!="add|change", GOTO="xorg_nvidia_end"
KERNEL!="nvidia", GOTO="xorg_nvidia_end"

ENV{xorg_driver}="nvidia", RUN+="/lib/udev/xorg-nvidia"

LABEL="xorg_nvidia_end"
