# network/if-up.d/kestrel
#!/bin/bash

[ "$IFACE" = lo ] && exit 0

[ "$MODE" != start ] && exit 0

[ -z "$KESTREL_SHARE" ] && KESTREL_SHARE=/usr/share/kestrel

. ${KESTREL_SHARE}/kestrel-lib

check_enabled || exit 0

[[ ${FRONTEND_IP} =~ ${ip_re} ]] &&
if ifconfig $IFACE | grep -Eq "addr:${FRONTEND_IP}\s"; then

    run_script "${KESTREL_SHARE}/services.d" > /dev/null
fi
