# network/if-up.d/kestrel
#! /bin/sh
# Reload the OpenSSH server when an interface comes up, to allow it to start
# listening on new addresses.

set -e

# Don't bother to restart sshd when lo is configured.
if [ "$IFACE" = lo ]; then
	exit 0
fi

# Only run from ifup.
if [ "$MODE" != start ]; then
	exit 0
fi

[ -z "$KESTREL_SHARE" ] && KESTREL_SHARE=/usr/share/kestrel

# Load kestrel_lib
. "$KESTREL_SHARE"/defaults-lib


if ifconfig $IFACE | grep -Eq "[^0-9]${FRONTEND_IP}[^0-9]"; then

    # idmapd sometimes fails to start on Ubuntu
    if [ -e /etc/init/idmapd.conf ]; then
        initctl status idmapd || initctl start idmapd
    fi

    [ -x /etc/init.d/dnsmasq ] &&
        invoke-rc.d dnsmasq restart

    invoke-rc.d kestrel_daemon restart

fi > /dev/null

exit 0
